# deploy/ircthing.service
# ircthing — self-hosted web IRC client.
#
# Expects: the binary at /usr/local/bin/ircd-web and config at
# /etc/ircthing/config.json (with "database": "/var/lib/ircthing/ircthing.db"
# — spelled out literally: systemd %-specifiers like %S only expand in unit
# directives, never inside the JSON config the service loads).
# No user setup needed — DynamicUser allocates one at runtime, and
# StateDirectory creates /var/lib/ircthing with correct ownership.
#
# The config is handed to the service as a systemd credential: it is
# copied into a private, service-only directory (mode 0400, owned by the
# runtime user) that no other service can read, and the app reads it from
# $CREDENTIALS_DIRECTORY/config.json. The on-disk /etc/ircthing/config.json
# stays owned by root and need not be world-readable.

[Unit]
Description=ircthing web IRC client
After=network-online.target
Wants=network-online.target

[Service]
DynamicUser=yes
StateDirectory=ircthing
WorkingDirectory=%S/ircthing

# systemd copies /etc/ircthing/config.json into $CREDENTIALS_DIRECTORY as
# "config.json"; the binary prefers that path when the env var is set.
LoadCredential=config.json:/etc/ircthing/config.json

# Keeps the Go heap inside the project's 72 MB RSS target. GOMEMLIMIT is a
# soft limit the runtime aims for; MemoryMax is the kernel's hard ceiling
# that actually contains the process (OOM-kills it) if something runs away.
# Headroom over the 72 MB target absorbs non-heap RSS and transient spikes.
Environment=GOMEMLIMIT=64MiB
MemoryMax=128M

ExecStart=/usr/local/bin/ircd-web
Restart=on-failure
RestartSec=5

# Hardening. Note DynamicUser=yes already implies NoNewPrivileges,
# PrivateTmp, ProtectSystem=strict, ProtectHome=read-only, RemoveIPC,
# and RestrictSUIDSGID — the lines below keep the stricter/explicit
# settings and everything DynamicUser does not cover.
ProtectHome=true
PrivateDevices=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# AF_NETLINK is required by the WireGuard egress dialer's UDP bind path (route
# lookup during socket setup); without it a WireGuard-proxied network fails to
# come up with EAFNOSUPPORT. Do NOT remove during hardening passes — the failure
# only manifests on WG-enabled networks, so it is easy to drop by mistake.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
